# launchd/com.gos.evolve-specs.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gos.evolve-specs</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/garyg/Documents/Documents - SG-LT674/Claude Working Folder/gOS/tools/evolve-claw.sh</string>
        <string>specs</string>
    </array>

    <!-- Mondays 9:13 (Weekday 1 = Monday in launchd) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key><integer>1</integer>
        <key>Hour</key><integer>9</integer>
        <key>Minute</key><integer>13</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/garyg/.claude/logs/evolve-specs-launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/garyg/.claude/logs/evolve-specs-launchd.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/garyg/.local/bin:/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/garyg</string>
    </dict>
</dict>
</plist>
